# common/etc/profile.d/aliases.sh
alias artisan="php '$APP_BASE_DIR'/artisan"
