# common/etc/s6-overlay/s6-rc.d/queue/run
#!/command/with-contenv sh

if [ "$QUEUE_ENABLED" = "true" ]; then
    echo "✅ Queue is enabled"
    queue-work
else
    echo "❌ Queue is disabled"
    s6-svc -Od .
    exit 0
fi
